# build/g-observer.service
[Unit]
Description=G-Observer Service
After=network.target

[Service]
User=root
ExecStartPre=/bin/sleep 5
ExecStart=/usr/local/bin/observer -config /usr/local/etc/observer/config.json
Restart=always

[Install]
WantedBy=multi-user.target
